# systemd/chungyak-collect.timer
[Unit]
Description=chungyak-alimi 배치 스케줄 (하루 2회: 08:00, 20:00)

[Timer]
OnCalendar=*-*-* 08,20:00:00
Persistent=true

[Install]
WantedBy=timers.target
